# portal/config/cron.d/portal
* *  * * *  root /opt/Portal-v3/portal/watchdog.sh
* *  * * *  root /opt/Portal-v3/portal/send.sh > /dev/null
